# launchd/com.goldfish.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.goldfish.heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/code/goldfish/launchd/goldfish-env.sh &amp;&amp; cd ~/code/goldfish &amp;&amp; hour=$(date +%H); if [ "$hour" -ge 9 ] &amp;&amp; [ "$hour" -le 22 ]; then pnpm run cli initiate -t heartbeat; else echo "Outside waking hours ($hour), skipping"; fi</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/tmp/goldfish-heartbeat.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/goldfish-heartbeat-err.log</string>
</dict>
</plist>
