# launchd/com.goldfish.daily-synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.goldfish.daily-synthesis</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source ~/code/goldfish/launchd/goldfish-env.sh &amp;&amp; ~/code/goldfish/scripts/daily-synthesis.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/goldfish-synthesis.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/goldfish-synthesis-err.log</string>
</dict>
</plist>
